# oresat-cand.service
[Unit]
Description=OreSat CAN daemon
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/oresat-cand can0
Restart=on-failure
User=root
Group=root

[Install]
WantedBy=multi-user.target
Alias=cand.service
