# debian/oresat-cand.prerm
#!/bin/sh -e
systemctl stop oresat-cand
systemctl disable oresat-cand
