# debian/oresat-cand.postinst
#!/bin/sh -e
systemctl daemon-reload
systemctl enable oresat-cand
systemctl start oresat-cand
