# debian/oresat-cand.postinst
#!/bin/sh -e

ORESAT_CONF_DIR=/etc/oresat
NODE_CONF=$ORESAT_CONF_DIR/node.conf
if [ ! -f "$NODE_CONF" ]; then
mkdir $ORESAT_CONF_DIR
cat >> $NODE_CONF << EOF
[Node]
NodeId=0x7C
CanInterface=can0
NetworkManager=false
EOF
endif

systemctl daemon-reload
systemctl enable oresat-cand
systemctl start oresat-cand
