# debian/oresat-canopend.postinst
#!/bin/sh -e
systemctl daemon-reload
systemctl enable oresat-canopend
systemctl start oresat-canopend
