# debian/oresat-canopend.prerm
#!/bin/sh -e
systemctl stop oresat-canopend
systemctl disable oresat-canopend
